# services/discovery-pipeline.service
[Unit]
Description=Polybot market-first discovery pipeline (discover/grade/plan/scan loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tstuv/poly/anthropic-exec-bot
EnvironmentFile=/home/tstuv/poly/anthropic-exec-bot/.env
Environment=PATH=/home/tstuv/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/tstuv/poly/anthropic-exec-bot/.venv/bin/python -m polybot.geopolitics run-discovery --config configs/discovery/geopolitics.yaml
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/tstuv/poly/anthropic-exec-bot/logs/discovery_pipeline.out
StandardError=append:/home/tstuv/poly/anthropic-exec-bot/logs/discovery_pipeline.out
SyslogIdentifier=discovery-pipeline

[Install]
WantedBy=default.target
