# services/geopolitics-fleet.service
[Unit]
Description=Polybot geopolitics fleet (discovery + one bot per eligible market)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tstuv/poly/anthropic-exec-bot
EnvironmentFile=/home/tstuv/poly/anthropic-exec-bot/.env
Environment=PATH=/home/tstuv/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Add --live only after a monitoring soak with fleet.position_mode alert_only.
ExecStart=/home/tstuv/poly/anthropic-exec-bot/.venv/bin/python -m polybot.geopolitics run-fleet --config configs/discovery/geopolitics.yaml
# Kill the whole control group so supervised bot subprocesses stop with the fleet.
KillMode=control-group
Restart=on-failure
RestartSec=30
StandardOutput=append:/home/tstuv/poly/anthropic-exec-bot/logs/geopolitics_fleet.out
StandardError=append:/home/tstuv/poly/anthropic-exec-bot/logs/geopolitics_fleet.out
SyslogIdentifier=geopolitics-fleet

[Install]
WantedBy=default.target
