# services/location-qatar.service
[Unit]
Description=Polybot Location Protection – Qatar Sept 30 YES position
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tstuv/poly/anthropic-exec-bot
EnvironmentFile=/home/tstuv/poly/anthropic-exec-bot/.env
Environment=PATH=/home/tstuv/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/tstuv/poly/anthropic-exec-bot/.venv/bin/python -m polybot.geopolitics run-location-protection --config configs/geopolitics/qatar-sept30-yes-protection.yaml --live
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/tstuv/poly/anthropic-exec-bot/logs/location_qatar_live.out
StandardError=append:/home/tstuv/poly/anthropic-exec-bot/logs/location_qatar_live.out
SyslogIdentifier=location-qatar

[Install]
WantedBy=default.target
