# deploy/polybot-fleet.service
# systemd unit for the fleet supervisor. Install:
#   sudo cp deploy/polybot-fleet.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now polybot-fleet
# A host reboot then restarts the whole fleet automatically; without this,
# a reboot silently stops all trading AND all position defense.
[Unit]
Description=Polybot geopolitics fleet supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust to the actual checkout path and user.
User=polybot
WorkingDirectory=/opt/anthropic-exec-bot
# Secrets (TELEGRAM_BOT_TOKEN, TELEGRAM_CHAT_ID, ANTHROPIC_API_KEY, wallet
# keys) live in an env file readable only by the service user -- never in
# the repo or the unit file.
EnvironmentFile=/etc/polybot/env
ExecStart=/opt/anthropic-exec-bot/.venv/bin/python -m polybot.geopolitics run-fleet --config configs/geopolitics/discovery.yaml --live
Restart=always
RestartSec=30
# The bots' own safety machinery (operator gate, process locks, journals)
# makes restart-always safe: a crashed supervisor resumes management, it
# never re-fires trades.

[Install]
WantedBy=multi-user.target
